# conf/supervisor/startup.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:apache2]
command=service apache2 start
redirect_stderr=true

[program:misc]
command=/etc/init.d/cron restart
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/cronlog.log
stdout_logfile=/var/log/cronlog.log

